# supervisord.conf

[supervisord]
#environment=LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8
#environment=LANG="ru_RU.utf8", LC_ALL="ru_RU.UTF-8", LC_LANG="ru_RU.UTF-8"

http_port=/var/tmp/supervisor.sock ; (default is to run a UNIX domain socket server)
;http_port=127.0.0.1:9001  ; (alternately, ip_address:port specifies AF_INET)
;sockchmod=0700              ; AF_UNIX socketmode (AF_INET ignore, default 0700)
;sockchown=nobody.nogroup     ; AF_UNIX socket uid.gid owner (AF_INET ignores)
;umask=022                   ; (process file creation umask;default 022)
logfile=/var/log/supervisor/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB       ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10          ; (num of main logfile rotation backups;default 10)
loglevel=info               ; (logging level;default info; others: debug,warn)
pidfile=/var/run/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false              ; (start in foreground if true;default false)
minfds=1024                 ; (min. avail startup file descriptors;default 1024)
minprocs=200                ; (min. avail process descriptors;default 200)

;nocleanup=true              ; (don't clean up tempfiles at start;default false)
;http_username=user          ; (default is no username (open system))
;http_password=123           ; (default is no password (open system))
;childlogdir=/tmp            ; ('AUTO' child log dir, default $TEMP)
;user=chrism                 ; (default is current user, required if root)
;directory=/tmp              ; (default is not to cd during start)
;environment=KEY=value       ; (key value pairs to add to environment)

[supervisorctl]
serverurl=unix:///var/tmp/supervisor.sock ; use a unix:// URL  for a unix socket
;serverurl=http://127.0.0.1:9001 ; use an http:// url to specify an inet socket
;username=chris              ; should be same as http_username if set
;password=123                ; should be same as http_password if set
;prompt=mysupervisor         ; cmd line prompt (default "supervisor")

; The below sample program section shows all possible program subsection values,
; create one or more 'real' program: sections to be able to control them under
; supervisor.



[program:uwsgi]
#environment=LANG="ru_RU.utf8", LC_ALL="ru_RU.UTF-8", LC_LANG="ru_RU.UTF-8"
command=/root/myblog/bin/uwsgi --ini /root/myblog/myblog/uwsgi.ini
stopsignal=QUIT
numprocs=1
#user=nginx
startretries =5
startsecs=5
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile = /var/log/uwsgi/log.log
stderr_logfile = /var/log/uwsgi/error.log


[program:celery]
command=/root/myblog/bin/celery  -A myblog -c 1 --beat worker -l info --workdir=/root/myblog/myblog --pidfile="/var/run/celery/%n.pid" --logfile="/var/log/celery/%n%I.log"
#command=/root/myblog/bin/celery multi start worker1 -A myblog -c 1 --beat worker -l info --workdir=/root/myblog/myblog --pidfile="/var/run/celery/%n.pid" --logfile="/var/log/celery/%n%I.log"
directory=/root/myblog/myblog
user=nginx
numprocs=1
startretries = 5
autostart=true
startsecs=10
autorestart=true
redirect_stderr=true
stopwaitsecs=60
stdout_logfile = /var/log/celery/celery.log
stderr_logfile = /var/log/celery/error.log
environment=C_FORCE_ROOT="true"
enviroment=PATH=/root/myblog/myblog
priority=998
#environment=LANG="ru_RU.utf8", LC_ALL="ru_RU.UTF-8", LC_LANG="ru_RU.UTF-8"

[program:redis]
command=/usr/local/bin/redis-server /root/myblog/myblog/redis.conf
user=nginx
numprocs=1
startretries = 5
autostart=true
startsecs=10
autorestart=true
#redirect_stderr=true
stopwaitsecs=20
#stdout_logfile = /var/log/celery/celery.log
#stderr_logfile = /var/log/celery/error.log






#[program:gunicorn]
#command=/root/myblog/bin/gunicorn main:application -c /root/myblog/gunicorn.conf.py
#directory=/root/myblog/myblog
#user=nginx
#startsecs=5
#autostart=true
#autorestart=true
#redirect_stderr=true
#command = /root/myblog/bin/gunicorn -w 2 --pythonpath=/root/myblog/bin,/root/myblog/myblog/blog --chdir=/root/myblog/myblog --reload --preload --error-logfile=/var/gunicorn/error.log --worker-connections 40000 --bind unix:/var/gunicorn/gunicorn.sock myblog.wsgi:application
#stdout_logfile = /var/log/ginicorn/gunicorn.log
#stderr_logfile = /var/log/gunicorn/error.log
#environment=LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8
#environment=LANG="ru_RU.utf8", LC_ALL="ru_RU.UTF-8", LC_LANG="ru_RU.UTF-8"
#enviroment=PATH=/root/myblog/myblog

#[program:jupyter]
#command=/root/myblog/bin/jupyter-notebook --no-browser --ip=0.0.0.0
#directory=/root/myblog/myblog
#user=root
#numprocs=1
#autostart=true
#startsecs=10
#autorestart=true
#redirect_stderr=true
#stopwaitsecs=600
##stdout_logfile = /var/log/celery/celery.log
#stderr_logfile = /var/log/celery/error.log
#environment=C_FORCE_ROOT="yes"
#enviroment=PATH=/root/myblog/myblog
#priority=9
